# infra/systemd/totolab-worker.service
[Unit]
Description=TotoLab worker — daily EPL match analysis via ai-debate bridge
After=network-online.target ai-debate.service
Wants=network-online.target ai-debate.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/toto-lab/worker
Environment=GOOGLE_APPLICATION_CREDENTIALS=/root/.firebase-sa.json
Environment=AI_DEBATE_URL=http://localhost:3000
Environment=AI_DEBATE_MODEL=claude-sonnet-4-6
ExecStart=/usr/bin/node /root/toto-lab/worker/runOnce.js
StandardOutput=journal
StandardError=journal
TimeoutStartSec=900
